# installer/systemd/homehub-server.service
[Unit]
Description=HomeHub local app, Google sync and setup portal
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=homehub
Group=homehub
WorkingDirectory=/opt/homehub/current
Environment=HOMEHUB_APP_ROOT=/opt/homehub/current
Environment=HOMEHUB_STATE_DIR=/var/lib/homehub
Environment=PYTHONPATH=/opt/homehub/current/backend
ExecStart=/opt/homehub/venv/bin/python -m homehub.app
Restart=always
RestartSec=3
# The service may invoke only the exact appliance helpers allowed by the
# HomeHub sudoers file. No shell or general systemctl access is granted.
NoNewPrivileges=false
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/homehub

[Install]
WantedBy=multi-user.target
